# backend/scripts/systemd/blackbow-backup.service
[Unit]
Description=BlackBow Associates Database Backup Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=newadmin
Group=newadmin
WorkingDirectory=/home/newadmin/projects/blackbow-associates/backend
Environment="NODE_ENV=production"
Environment="BACKUP_DIR=/home/newadmin/projects/blackbow-associates/backend/backups"
Environment="LOG_FILE=/var/log/desaas/blackbow-backup.log"

# Load environment variables from .env file
EnvironmentFile=-/home/newadmin/projects/blackbow-associates/backend/.env

# Execute backup script
ExecStart=/home/newadmin/projects/blackbow-associates/backend/scripts/backup.sh

# Standard output and error go to journal
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/newadmin/projects/blackbow-associates/backend/backups /var/log/desaas

# Restart settings (only on failure, not for oneshot)
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
